# root/etc/services.d/noip-rfc2136/run
#!/usr/bin/with-contenv bash

exec \
        s6-setuidgid abc python3 /opt/venv/noip-rfc2136.py
