# root/etc/services.d/noip-rfc2136/run
#!/usr/bin/with-contenv bash

cd /opt/venv/
exec \
        s6-setuidgid abc python3 noip-rfc2136.py
